# config/includes.chroot/etc/systemd/system/waydroid-binder-setup.service
[Unit]
Description=Prepare Android binder devices for Waydroid
DefaultDependencies=no
Before=waydroid-container.service
After=systemd-modules-load.service local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/waydroid-binder-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
